# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

#nc -d 192.168.11.34 30005 | nc 127.0.0.1 30104 &
#nc -d 192.168.11.32 30005 | nc 127.0.0.1 30104 &
#nc -d 192.168.11.44 30005 | nc 127.0.0.1 30104 &
#nc -d 192.168.11.8  30005 | nc 127.0.0.1 30104 &
#nc -d 192.168.11.8  30005 | nc 104.155.90.204 30104 &
/home/pi/connect.dump.sh

exit 0
